# base/Linux/ROLES/squid/etc/network/if-up.d/openntpd
#!/bin/sh

if [ "${METHOD}" = "loopback" ]
then
	exit 0
fi

invoke-rc.d openntpd force-reload
